# nasa-apod.service
[Unit]
Description=NASA APOD Downloader

[Service]
ExecStart=/home/coconnor/bin/nasa-apod.sh
