# systemd/user/gnome-session-manager@.service.d/overrides.conf
[Service]
ExecStart=${gnome-session-binary} --systemd-service --session=$i --debug
