# pkgs/systemd/user/gnome-session-manager@.service.d/overrides.conf
[Service]
ExecStart=${gnome-session-binary} --session=%i --debug

